# etc/udev/rules.d/90-obsbot-tiny3.rules
ACTION=="add", SUBSYSTEM=="video4linux", ATTRS{idVendor}=="3564", ATTRS{idProduct}=="ff02", RUN+="/usr/local/bin/obsbot-fix-whitebalance.sh"
