# home/.config/autostart/xmousepasteblock.desktop
[Desktop Entry]
Type=Application
Name=xmousepasteblock
Exec=xmousepasteblock
X-GNOME-Autostart-enabled=true
